# namechanger.service
[Unit]
Description=namechanger

[Service]
Restart=on-failure
ExecStart=/usr/bin/namechanger -c /etc/namechanger.yml